# ganetimgr-alpine/assets/etc/supervisor.d/uwsgi.ini
[program:uwsgi]
user=ganetimgr
stopsignal=INT
command=uwsgi --ini uwsgi.ini
