# ganetimgr-alpine/assets/etc/supervisor.d/watcher.ini
[program:watcher]
user=ganetimgr
command=/srv/ganetimgr/watcher.py -p /tmp/ganetimgr-watcher.pid -l /tmp/ganetimgr-watcher.log -f
