# assets/etc/supervisor.d/nginx.ini
[program:nginx]
user=root
directory=/var
command=/usr/sbin/nginx
