# assets/etc/supervisor.d/uwsgi.ini
[program:uwsgi]
user=ganetimgr
command=uwsgi --http11-socket :8080 --plugin=python3 --wsgi ganetimgr.wsgi:application --tcp-nodelay --so-keepalive --die-on-term
